# kms/deploy/mx93/kms-api.service
[Unit]
Description=AirAccount KMS API Server
# dirf-repair.service must run first: it stops/starts tee-supplicant to repair
# a corrupt dirf.db. Running that logic in ExecStartPre here would deadlock because
# systemd holds the tee-supplicant unit lock while starting kms-api.service, and
# restarting tee-supplicant from ExecStartPre blocks on a transaction that never
# resolves. The separate oneshot unit avoids this dependency cycle entirely.
After=network.target tee-supplicant@teepriv0.service dirf-repair.service
Requires=tee-supplicant@teepriv0.service dirf-repair.service

[Service]
Environment=KMS_ENABLE_CREATEKEY=1
# Operator config (optional). Drop a kms.env here to set your own WebAuthn
# relying party (KMS_RP_ID / KMS_RP_NAME / KMS_ORIGIN), rate limits, etc. — no
# recompile needed. API-key auth is fail-closed by default (a key is required);
# set KMS_ALLOW_OPEN_MODE=1 only for an intentionally unauthenticated dev box.
# The leading
# `-` makes it optional: if the file is absent the service uses built-in
# defaults (rpId=aastar.io). See deploy/mx93/kms.env.example for the template.
EnvironmentFile=-/etc/airaccount/kms.env
Type=simple
User=root
WorkingDirectory=/root/AirAccount
ExecStart=/root/AirAccount/target/release/kms-api-server
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/kms-api.log
StandardError=append:/var/log/kms-api.log
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

[Install]
WantedBy=multi-user.target
